# services/jotigames-frontend.service
[Unit]
Description=JotiGames Frontend (Vite Preview)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/jotigames.nl/frontend
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStartPre=/bin/bash -lc 'node -e "const v=parseInt(process.versions.node.split(\".\")[0],10); if (!Number.isInteger(v) || v < 20) process.exit(1);"'
ExecStartPre=/bin/bash -lc 'test -d /var/www/jotigames.nl/frontend/dist'
ExecStart=/usr/bin/env node /var/www/jotigames.nl/frontend/node_modules/vite/bin/vite.js preview --host 0.0.0.0 --port 4173
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
