# services/jotigames-frontend.service
[Unit]
Description=JotiGames Frontend (Vite Preview)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/jotigames.nl/frontend
ExecStartPre=/bin/bash -lc '/usr/bin/env node --input-type=module -e "await Promise.resolve()"; /usr/bin/env node -e "const m=process.versions.node.match(/^([0-9]+)\.([0-9]+)/); const major=Number(m?.[1]||0); const minor=Number(m?.[2]||0); if (!Number.isInteger(major) || !Number.isInteger(minor) || major < 24 || (major === 24 && minor < 14)) process.exit(1);"'
ExecStartPre=/bin/bash -lc 'test -d /var/www/jotigames.nl/frontend/dist'
ExecStart=/bin/bash -lc 'exec /usr/bin/env node /var/www/jotigames.nl/frontend/node_modules/vite/bin/vite.js preview --host 0.0.0.0 --port 4173'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
