# services/jotigames-socketserver.service
[Unit]
Description=JotiGames WebSocket Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/jotigames.nl/ws
Environment=NODE_ENV=production
Environment=NODE_BIN=/root/.nvm/versions/node/v24.14.0/bin/node
Environment=PATH=/root/.nvm/versions/node/v24.14.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStartPre=/bin/bash -lc 'NODE_BIN="${NODE_BIN:-$(command -v node)}"; test -x "$NODE_BIN"; "$NODE_BIN" --input-type=module -e "await Promise.resolve()"; "$NODE_BIN" -e "const m=process.versions.node.match(/^([0-9]+)\.([0-9]+)/); const major=Number(m?.[1]||0); const minor=Number(m?.[2]||0); if (!Number.isInteger(major) || !Number.isInteger(minor) || major < 24 || (major === 24 && minor < 14)) process.exit(1);"'
ExecStart=/bin/bash -lc 'NODE_BIN="${NODE_BIN:-$(command -v node)}"; exec "$NODE_BIN" /var/www/jotigames.nl/ws/server.js'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target

